# deploy/com.ufii.folia-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ufii.folia-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/ufi/dev/folia-major/server/src/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ufi/dev/folia-major/server</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>ENABLE_GENERAL_UNBLOCK</key>
        <string>true</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ufi/Library/Logs/folia-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ufi/Library/Logs/folia-server.log</string>
</dict>
</plist>
